# daraz-project/systemd/system/ngrok-tunnel.service
[Unit]
Description=Ngrok Tunnel for Daraz API
After=network.target daraz-mcp.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/daraz-scraper/daraz-project
ExecStart=/bin/bash /home/ubuntu/daraz-scraper/daraz-project/start_ngrok.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target